# chef-repo/site-cookbooks/factorialapp/files/default/factorial-api.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec python /opt/devops-project/factorial.py
